# src/metadata/entries/autostart/com.deepin.music.autostart.desktop
[Desktop Entry]
Type=Application
Exec=/opt/apps/com.deepin.demo/files/bin/com.deepin.demo
GenericName=App Store Demo
Icon=com.deepin.demo